# backend/crontab
*/15 * * * * root cd /app/backend && uv run python scripts/rollup_trending_views.py >> /proc/1/fd/1 2>> /proc/1/fd/2
