# backend/crontab
*/15 * * * * root cd /app/backend && uv run python scripts/rollup_trending_views.py >> /proc/1/fd/1 2>> /proc/1/fd/2
*/15 * * * * root cd /app/backend && uv run python scripts/rollup_download_stats.py >> /proc/1/fd/1 2>> /proc/1/fd/2
* * * * * root cd /app/backend && uv run python scripts/upload_links_to_servers.py >> /proc/1/fd/1 2>> /proc/1/fd/2
